# task2/src/supervisord.conf
# supervisord 进程管理 —— 管控 collector（采集）+ janitor（清理）常驻
# 日志统一转发到容器 stdout/stderr（docker logs 可见）

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

# C1 采集器：perf record 持续采集 + 60s 轮转
[program:collector]
command=python3 -m profiler.collector
autorestart=true
startretries=5
startsecs=3
# 输出转发到 supervisord 的 stdout（容器 stdout）
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

# C2 清理器：定时删过期采样文件
[program:janitor]
command=python3 -m profiler.janitor
autorestart=true
startsecs=3
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
